# layout/DEBIAN/postinst
#!/bin/bash

echo "Restarting/loading daemon..."
killall ambid || true
chown 0:0 /Library/LaunchDaemons/me.k3a.ambid.plist
launchctl load -w /Library/LaunchDaemons/me.k3a.ambid.plist || true
